# etc/system/met_moai_update.timer
[Unit]
Description=Run met_moai_update.service every hour

[Timer]
OnCalendar=*-*-* *:19:00
